# rc.local

sudo nano /etc/rc.local
---

service fail2ban start
sudo /etc/init.d/shadowsocks start
sudo service v2ray start
sudo systemctl restart noip2
exit 0

---
